# ansible/roles/db/templates/clustercheck@.service.j2
[Unit]
Description=Cliustercheck

[Service]
ExecStart=-/usr/local/bin/clustercheck clustercheck "{{ clustercheck_password }}"
StandardInput=socket
