# ansible/roles/db/files/clustercheck.socket
[Unit]
Description=Galera cluster check

[Socket]
ListenStream=9200
Accept=yes

[Install]
WantedBy=sockets.target
